# conf/anote.conf
[program:anote]
command=/root/anote-node
process_name=%(program_name)s
directory=/root
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/anote.log